# agent/scripts/da-invent-agent.service
[Unit]
Description=DA-INVENT Remote Agent
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=da-invent-agent
Group=da-invent-agent
WorkingDirectory=/opt/da-invent-agent
Environment=DA_INVENT_AGENT_CONFIG=/etc/da-invent-agent/config.yml
ExecStart=/opt/da-invent-agent/venv/bin/da-invent-agent
Restart=on-failure
RestartSec=5

# Hardening
ProtectSystem=strict
ReadWritePaths=/var/log/da-invent-agent
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
AmbientCapabilities=CAP_NET_RAW CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_BIND_SERVICE
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
