# deploy/da-invent.service
[Unit]
Description=DA-INVENT - IP Address Management
After=network.target

[Service]
Type=simple
# Nel container (LXC/VM dedicata) l’app gira come root così nmap può usare
# scansioni UDP (-sU) e ICMP senza capability extra. Per installazioni bare metal
# con utente dedicato impostare DA_INVENT_SERVICE_USER=da-invent in install.sh
# e concedere CAP_NET_RAW / CAP_NET_ADMIN (vedi README).
User=root
Group=root
WorkingDirectory=/opt/da-invent
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/opt/da-invent/.env.local
ExecStart=/opt/da-invent/node_modules/.bin/tsx /opt/da-invent/server.ts
Restart=on-failure
RestartSec=5
# No PrivateTmp: evita 226/NAMESPACE su LXC/Proxmox (vedi scripts/install.sh).

[Install]
WantedBy=multi-user.target
