# deploy/da-invent-update.timer
[Unit]
Description=DA-INVENT auto-update giornaliero

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=900
Persistent=true
Unit=da-invent-update.service

[Install]
WantedBy=timers.target
